# deploy/site-cookbooks/webapp/templates/default/supervisior.conf.erb
[program:backend]
command = <%= node.default['webapp']['gunicorn_script']%>                     ; Command to start app
user = <%= node[:user]%>                                               ; User to run as
autostart = true
autorestart = true
redirect_stderr = true                                      ; Save stderr in the same log
